# dot_config/systemd/user/swww.service
[Unit]
Description=Wallpaper daemon (swww) — user-session
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/swww-daemon
ExecStartPost=/usr/bin/bash -c '/usr/bin/swww img --transition-type random "$(find "$HOME/.config/background" -type f | shuf -n1)"'
Restart=on-failure

[Install]
WantedBy=default.target
